# docker/prod/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
# user=www-data
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
# pidfile=/tmp/supervisord.pid

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm/access.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/php-fpm/error.log
stderr_logfile_maxbytes=0
exitcodes=0
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/access.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/nginx/error.log
stderr_logfile_maxbytes=0
exitcodes=0
stopsignal=QUIT

# [program:certbot]
# process_name=certbot
# command=sh -c "/usr/local/bin/certbot.sh"
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/certbot.log
# stdout_logfile_maxbytes=0
# stderr_logfile=/var/log/certbot_error.log
# stderr_logfile_maxbytes=0
# exitcodes=0
# stopsignal=QUIT

[program:horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/artisan horizon
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/horizon.log
stopwaitsecs=3600

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
